# etc/skel/systemd/system/disable-touchpad.service
[Unit]
Description=Alway disable touchpad

[Service]
ExecStart=/usr/local/bin/disable-touchpad
Type=simple

[Install]
WantedBy=default.target
